# roles/docker-container-scanner/files/crontab
@reboot root    cd /usr/local/docker && docker stack deploy traefik --compose-file=docker-compose.yml
@reboot root    /usr/local/sbin/scanner.sh 
