# roles/docker-container-scanner/files/crontab
@reboot root    docker stack deploy traefik --compose-file=/usr/local/docker/docker-compose.yml
@reboot root    /usr/local/sbin/scanner.sh 
