# roles/docker-container-scanner/files/crontab
@reboot root    /usr/local/sbin/scanner.sh 
