# roles/traefik/files/crontab
@reboot root    docker stack deploy traefik --compose-file=/usr/local/docker/docker-compose-traefik.yml
